# profile.d/watchman.sh
#!/usr/bin/env bash

set -eu

export PATH="$PATH:$HOME/vendor"
